# rpi_ssh/create.service
[Unit]
Description=create
After=multi-user.target
[Service]
Type=idle
ExecStart=/root/create.sh
[Install]
WantedBy=multi-user.target
